# homework18/files/borg-backup.timer
[Unit]
Description=Borg Backup
Requires=borg-backup.service

[Timer]
Unit=borg-backup.service
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
